# support/dev/redbox.cron
#!/bin/sh
sudo -Hu redbox /home/redbox/deploy.sh redbox
sudo -Hu redbox /home/redbox/deploy.sh mint
